# cloudflared.service
[Unit]
Description=Cloudflared tunnel — ISP Monitor (URL tetap)
After=network.target

[Service]
Type=simple
User=bahcron
EnvironmentFile=/home/bahcron/isp-monitor/.cf-tunnel-token
ExecStartPre=/bin/sh -c 'until curl -sf http://localhost:8000/healthz >/dev/null; do sleep 2; done'
ExecStart=/home/bahcron/bin/cloudflared --no-autoupdate tunnel run --token ${CF_TUNNEL_TOKEN}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target