# isp-monitor-probe.service
[Unit]
Description=ISP Monitor — Probe Region (multi-region)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/isp-monitor
Environment=NODE_ENV=production
EnvironmentFile=/opt/isp-monitor/.env
ExecStart=/usr/bin/node src/probe.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
